# debian/postrm
#! /bin/sh
set -e

case "$1" in
purge )
  update-rc.d p4 remove >/dev/null
  ;;
esac
